# alternc-httpperm.cron


* * * * * root cd /usr/lib/alternc && php httpperm_cron.php &>/dev/null


